# rootfs/etc/s6-overlay/s6-rc.d/mountpoint-s3/run
#!/command/execlineb -P

with-s3-env
importas -i BUCKET MOUNTPOINT_BUCKET
importas -i DIR MOUNTPOINT_DIR
s6-notifyoncheck -n 0 -c "mountpoint -q ${DIR}"
mount-s3 \
    --foreground \
    --force-path-style \
    --upload-checksums off \
    --cache /var/cache/mount-s3 \
    $BUCKET $DIR

    # --debug \
    # --debug-crt \