# docker/s6-rootfs/etc/s6-overlay/s6-rc.d/nginx-conf/up
# code: language=shellscript

with-contenv
redirfd -w 1 /etc/nginx/nginx.conf
redirfd -r 0 /etc/nginx/nginx.conf.template
envsubst '$AWS_UPSTREAM_HOST' '$AWS_UPSTREAM_PROTOCOL'
