# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P

with-s3-env
s6-notifyoncheck -n 0 -c "nc -z localhost 54321"
nginx -g "daemon off;"
