# rootfs/etc/s6-overlay/s6-rc.d/pre-mountpoint-s3-1/up
# code: language=shellscript

with-s3-env
importas -i DIR MOUNTPOINT_DIR_1
if { mkdir -p /var/cache/mount-s3-1 }
if { mkdir -p "${DIR}" }
mount-s3 --version
