# rootfs/etc/s6-overlay/s6-rc.d/nginx-conf/up
# code: language=shellscript

with-s3-env
redirfd -w 1 /etc/nginx/nginx.conf
redirfd -r 0 /etc/nginx/nginx.conf.template
envsubst "$AWS_UPSTREAM_SCHEME,$AWS_UPSTREAM_HOST,$AWS_UPSTREAM_PORT"
