# docker/s6-rootfs/etc/s6-overlay/s6-rc.d/pre-mountpoint-s3/up
# code: language=shellscript

with-contenv
if { mkdir -p /var/cache/mount-s3 }
if { mkdir -p /mnt/pub }
mount-s3 --version
